# salt/redmine-backup.service
[Unit]
After=postgresql.service

[Service]
User=www-data
ExecStart=/opt/redmine-3.3.2/backup
Type=oneshot
