# salt/tt-rss-backup.service
[Unit]
After=postgresql.service

[Service]
User=www-data
ExecStart=/opt/tt-rss/backup
Type=oneshot
